# config/production/supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock

[supervisord]
logfile=/var/log/crystalcog/supervisord.log
logfile_maxbytes=50MB
logfile_backups=10
loglevel=info
pidfile=/tmp/supervisord.pid
nodaemon=true
minfds=1024
minprocs=200
user=crystalcog

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[program:crystalcog-api]
command=/usr/local/bin/demo
directory=/app
user=crystalcog
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/crystalcog/api.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
environment=CRYSTAL_ENV=production,LOG_LEVEL=info

[program:cogserver]
command=/usr/local/bin/demo_cogserver
directory=/app
user=crystalcog
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/crystalcog/cogserver.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
environment=CRYSTAL_ENV=production,COGSERVER_PORT=17001

[program:ai-integration]
command=/usr/local/bin/demo_ai_integration
directory=/app
user=crystalcog
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/crystalcog/ai-integration.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
environment=CRYSTAL_ENV=production

[program:python-kernel]
command=python3 python_cognitive_kernel.py
directory=/app
user=crystalcog
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/crystalcog/python-kernel.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
environment=PYTHONPATH=/home/crystalcog/.local/lib/python3.11/site-packages

[program:flask-app]
command=python3 app.py
directory=/app
user=crystalcog
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/crystalcog/flask-app.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
environment=FLASK_ENV=production,PYTHONPATH=/home/crystalcog/.local/lib/python3.11/site-packages

[group:crystalcog-services]
programs=crystalcog-api,cogserver,ai-integration,python-kernel,flask-app